# ansible/roles/nodejs/tasks/templates/nodeapp.service.j2
[Unit]
Description=Node.js Application
After=network.target

[Service]
Type=simple
User={{ app_user }}
Group={{ app_group }}
WorkingDirectory=/opt/{{ app_name }}
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target